# deploy/dunamismax-site.service
[Unit]
Description=dunamismax.com Spring Boot site
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=sawyer
Group=sawyer
WorkingDirectory=/opt/dunamismax-site
Environment="JAVA_HOME=/usr/lib/jvm/temurin-21-jdk-amd64"
Environment="JAVA_TOOL_OPTIONS=-XX:MaxRAMPercentage=70 -XX:+UseG1GC -XX:+ExitOnOutOfMemoryError"
EnvironmentFile=-/opt/dunamismax-site/site.env
ExecStart=/usr/lib/jvm/temurin-21-jdk-amd64/bin/java -jar /opt/dunamismax-site/dunamismax-site.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dunamismax-site

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/dunamismax-site
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
