# deploy/dunamismax-site.service
[Unit]
Description=dunamismax.com Rust site
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=dunamismax
Group=dunamismax
WorkingDirectory=/opt/dunamismax-site
EnvironmentFile=-/opt/dunamismax-site/site.env
ExecStart=/opt/dunamismax-site/dunamismax-site
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dunamismax-site

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dunamismax-site
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
